# dev-util/gitlab-runner/files/gitlab-runner
#!/sbin/openrc-run

# to add to startup simply run rc-update add gitlab-runner default
# this script assumes you installed the runner in /usr/bin.
NAME="gitlab-runner"
PIDFILE="/var/run/$NAME.pid"
DAEMON="/usr/bin/$NAME"
LOGFILE="/var/log/$NAME.log"

USER="$NAME"
GROUP="$NAME"

depend() {
        after docker
        need net
        need localmount
}

start() {
        ebegin "Starting $NAME"
        start-stop-daemon --start --exec /usr/bin/$NAME --pidfile "$PIDFILE" --chdir "/home/$USER" --background --make-pidfile --exec "$DAEMON" -- run &>> $LOGFILE
        eend $?
}

stop() {
        ebegin "Stopping $NAME"
        start-stop-daemon --stop --pidfile "$PIDFILE" --user "$USER" --exec "$DAEMON" --quiet
        eend $?
}
